# deploy/systemd/hedgents-orchestrator.service
[Unit]
Description=Hedgents orchestrator-daemon (regime-aware allocator)
After=network-online.target hedgents-dashboard.service
Wants=network-online.target
Requires=hedgents-dashboard.service
PartOf=hedgents.target

[Service]
StandardOutput=append:/var/lib/hedgents/logs/orchestrator.log
StandardError=append:/var/lib/hedgents/logs/orchestrator.log
Type=simple
User=hedgents
Group=hedgents
WorkingDirectory=/var/lib/hedgents
EnvironmentFile=/etc/hedgents/hedgents.env
# Optional operator-supplied flags. When `/etc/hedgents/orchestrator.conf`
# defines `EXECUTE_FLAGS=...`, that string is substituted into the command
# line below — typically:
#
#     EXECUTE_FLAGS=--execute --targets-json=/etc/hedgents/orchestrator-targets-mainnet.json \
#         --cooldown-secs=300 --max-action-fraction=0.50 --min-action-usd=5.0 --stale-slack=1.10
#
# Without the conf file (default), `EXECUTE_FLAGS` is empty and the daemon
# runs in dry-run / observe-only mode. The `-` prefix on EnvironmentFile
# means systemd does NOT fail to start if the conf file is missing.
#
# Persisting these flags in /etc/hedgents/orchestrator.conf survives
# `install-hedgents.sh` re-runs, which overwrite this unit file from the
# release tarball but never touch /etc/hedgents/.
#
# Allocator v2 drift mode (rc22+):
# To switch from rc21 greedy allocation to operator-set target weights,
# append `--target-weights="<spec>"` to EXECUTE_FLAGS. Example:
#
#     EXECUTE_FLAGS=--execute --targets-json=... --cooldown-secs=300 \
#         --max-action-fraction=0.50 --min-action-usd=5.0 --stale-slack=1.10 \
#         --target-weights="stable_yield=0.30,multiply=0.30,hedgedjlp=0.40" \
#         --min-drift-bps=200
#
# Empty/unset `--target-weights` keeps greedy mode. The orchestrator
# refuses to boot on a malformed spec rather than silently falling back
# to greedy — a typo in the conf shouldn't quietly disable the feature
# you thought you enabled.
EnvironmentFile=-/etc/hedgents/orchestrator.conf
Environment=RUST_LOG_FORMAT=json
ExecStart=/opt/hedgents/bin/orchestrator-daemon \
    --secrets-dir=/var/lib/hedgents/secrets \
    --listen=/ip4/0.0.0.0/tcp/19317 \
    --bootstrap=/ip4/127.0.0.1/tcp/19312 \
    --beacon-interval-secs=10 \
    --api-base=http://127.0.0.1:7700 \
    --tick-interval-secs=60 \
    --audit-log=/var/lib/hedgents/logs/orchestrator-audit.jsonl \
    $EXECUTE_FLAGS
Restart=always
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=120

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hedgents
PrivateTmp=true

[Install]
WantedBy=hedgents.target
