# deploy/systemd/hedgents-stable-yield-live.service
[Unit]
Description=Hedgents stable-yield-daemon (LIVE mainnet — Kamino USDC supply)
After=network-online.target hedgents-multiply-live.service
Wants=network-online.target
PartOf=hedgents-live.target
Conflicts=hedgents-stable-yield.service

[Service]
StandardOutput=append:/var/lib/hedgents/logs/stable-yield-live.log
StandardError=append:/var/lib/hedgents/logs/stable-yield-live.log
Type=simple
User=hedgents
Group=hedgents
WorkingDirectory=/var/lib/hedgents
EnvironmentFile=/etc/hedgents/hedgents.env
ExecStart=/opt/hedgents/bin/stable-yield-daemon \
    --secrets-dir=/var/lib/hedgents/secrets \
    --rpc-url=${RPC_URL} \
    --listen=/ip4/0.0.0.0/tcp/19313 \
    --bootstrap=/ip4/127.0.0.1/tcp/19312 \
    --network=mainnet \
    --i-understand-this-is-mainnet \
    --orchestrator-agent-id=${ORCHESTRATOR_PUBKEY} \
    --beacon-interval-secs=5 \
    --simulate-only=false \
    --require-approval=true \
    --telemetry-log=/var/lib/hedgents/logs/stable-yield-live-pnl.jsonl \
    --telemetry-market=7u3HeHxYDLhnCoErrtycNokbQYbWGzLs6JSDqGAv5PfF \
    --telemetry-interval-secs=30
Restart=on-failure
RestartSec=5s
StartLimitBurst=10
StartLimitIntervalSec=120

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hedgents
PrivateTmp=true

[Install]
WantedBy=hedgents-live.target
