# deploy/systemd/hedgents.target
[Unit]
Description=Hedgents fleet (six daemons + dashboard + paper-trade loop)
Wants=hedgents-multiply.service hedgents-stable-yield.service hedgents-hedgedjlp.service hedgents-riskwatcher.service hedgents-researcher.service hedgents-orchestrator.service hedgents-dashboard.service hedgents-paper-trade.service
After=hedgents-multiply.service hedgents-stable-yield.service hedgents-hedgedjlp.service hedgents-riskwatcher.service hedgents-researcher.service hedgents-orchestrator.service hedgents-dashboard.service hedgents-paper-trade.service

[Install]
WantedBy=multi-user.target
